# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
nodaemon = true

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; KYPO specific configuration

[program:redis]
command=redis-server /etc/redis.conf
priority=1
stopsignal=TERM
autostart=true
autorestart=true

[program:postgresql]
command=/usr/bin/postgres
user=postgres
priority=1
stopsignal=TERM
autostart=true
autorestart=true

[program:sandbox-service]
command=bin/run-sandbox-service.sh
environment=LC_ALL=C.UTF-8,LANG=C.UTF-8
priority=2
directory=/app
stopsignal=TERM
autostart=true
autorestart=true

[program:sandbox-service-rq-worker-openstack]
command=pipenv run python3 manage.py rqworker openstack
process_name=%(program_name)s-%(process_num)s
environment=LC_ALL=C.UTF-8,LANG=C.UTF-8
numprocs=2
directory=/app
stopsignal=TERM
autostart=true
autorestart=true

[program:sandbox-service-rq-worker-ansible]
command=pipenv run python3 manage.py rqworker ansible
process_name=%(program_name)s-%(process_num)s
environment=LC_ALL=C.UTF-8,LANG=C.UTF-8
numprocs=2
directory=/app
stopsignal=TERM
autostart=true
autorestart=true

[program:sandbox-service-rq-worker-default]
command=pipenv run python3 manage.py rqworker default
process_name=%(program_name)s-%(process_num)s
environment=LC_ALL=C.UTF-8,LANG=C.UTF-8
numprocs=2
directory=/app
stopsignal=TERM
autostart=true
autorestart=true

[group:sandbox-service]
programs=redis,postgresql,sandbox-service,sandbox-service-rq-worker-openstack,sandbox-service-rq-worker-ansible,sandbox-service-rq-worker-default

# vim: set syn=dosini :
